# piperpal.cron
0 * * * * python /home/ole/bin/piperpal
